# limit-boost.service
[Unit]
Description=Limit CPU boost based on cpu temps

[Service]
ExecStart=/usr/local/bin/limit-boost

[Install]
WantedBy=multi-user.target
