# rchill.service
[Unit]
Description=RyzenChill - Limit CPU boost based on cpu temps

[Service]
ExecStart=/usr/local/bin/rchill

[Install]
WantedBy=multi-user.target